# ansible/files/vpslite-backup.service
[Unit]
Description="Automatic data backup"

[Service]
Type=oneshot
ExecStart=/srv/backup/backup.sh
WorkingDirectory=/srv/backup
